# config/dht_sensor_exporter_nodejs.service
[Unit]
Description=Nodejs Prometheus exporter for DHT sensor information

[Service]
User=dht_sensor_exporter
Type=simple
ExecStart=node /opt/dht_sensor_exporter_nodejs/index.js

[Install]
WantedBy=multi-user.target
